# cgproxy.service
[Unit]
Description=cgproxy service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/cgproxyd

[Install]
WantedBy=multi-user.target
