# control/postinst
#!/bin/sh
systemctl enable --now cgproxy.service
