# control/prerm
#!/bin/sh
systemctl disable --now cgproxy.service 
